# cron/jobs.crontab
*/15 * * * * /home/danny/projects/dotfiles/cron/aws-azure-login.sh
