# cron/jobs.crontab
*/15 * * * * ~/projects/dotfiles/cron/aws-azure-login.sh

